# meta-digi-arm/recipes-kernel/kernel-module-qualcomm/kernel-module-qualcomm/80-sdio-qcom.rules
# Load Qualcomm wireless module
SUBSYSTEM=="sdio", ACTION=="add", ENV{MODALIAS}=="sdio:c00v0271d050A" RUN+="/etc/udev/scripts/qualcomm.sh"
