# meta-digi-arm/recipes-bsp/firmware-murata-infineon/firmware-murata-infineon/cyw55512-bluetooth.service
[Unit]
Description=CYW55512 Bluetooth Initialization
After=network.target suspend.target

[Service]
Type=forking
ExecStart=/etc/cyw55512-bluetooth

[Install]
WantedBy=multi-user.target suspend.target
