# meta-digi-dey/recipes-connectivity/modemmanager/modemmanager/late-modems-scan.timer
[Unit]
Description=Late modems scan timer

[Timer]
OnActiveSec=40s

[Install]
WantedBy=multi-user.target
