# meta-digi-dey/recipes-multimedia/pulseaudio/pulseaudio/pulseaudio-system.service
[Unit]
Description=System-wide sound service
After=dbus.service

[Service]
Type=forking
PIDFile=/run/pulse.pid
ExecStart=/etc/pulseaudio-init
ExecStop=kill -HUP $MAINPID

[Install]
WantedBy=multi-user.target
