# meta-digi-dey/recipes-connectivity/bluez/bluez5-5.41/bluetooth-init.service
[Unit]
Description=Initialization of the bluetooth chip

[Service]
Type=forking
ExecStart=/etc/bluetooth-init start
ExecStop=/etc/bluetooth-init stop

[Install]
RequiredBy=bluetooth.service
