# meta-digi-dey/recipes-core/busybox/busybox/busybox-httpd.service.in
[Unit]
Description=Busybox HTTP daemon
After=network.target

[Service]
Type=forking
ExecStart=@sbindir@/httpd -h /srv/www
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
