# meta-digi-dey/dynamic-layers/meta-ml/recipes-libraries/eiq-examples/files/service/eiqdemo.service
[Unit]
Description=Initialization of EiQ DMS demo
After=multi-user.target

[Service]
Type=forking
ExecStart=/etc/demos/scripts/launch_eiq_demo.sh

[Install]
WantedBy=multi-user.target
