# meta-digi-dey/recipes-connectivity/bluez/bluez5-init/bluetooth-init.service
[Unit]
Description=Initialization of the bluetooth chip
ConditionPathIsDirectory=/proc/device-tree/bluetooth
Before=bluetooth.service

[Service]
Type=forking
ExecStart=/etc/bluetooth-init start
ExecStop=/etc/bluetooth-init stop

[Install]
WantedBy=bluetooth.service
